# apps/runner-agent/deployment/axion-agent.service
[Unit]
Description=Axion Runner Agent
Documentation=https://github.com/axion/axion-stack
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/axion-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=axion-agent

# Environment
Environment="RUST_LOG=info"
EnvironmentFile=-/etc/axion/agent.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/axion /var/log/axion
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryLimit=512M
CPUQuota=200%

[Install]
WantedBy=multi-user.target
